# docker_cron/cron.d/crontab
* * * * * docker kill -s HUP 4f3e3234299984474993ad2ae177e7cc8521280d60e9ee3f0ff378ed0a63c5c2
